# deploy/mac/com.openfrontier.cloudflared.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openfrontier.cloudflared</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/YOUR_USER/Documents/OpenFrontier/deploy/mac/start-open-frontier-cloudflared.zsh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/YOUR_USER/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OPEN_FRONTIER_CLOUDFLARED_CONFIG</key>
    <string>/Users/YOUR_USER/.cloudflared/open-frontier-beta.yml</string>
    <key>OPEN_FRONTIER_CLOUDFLARE_TUNNEL_NAME</key>
    <string>open-frontier-beta</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/Library/Logs/open-frontier-cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/Library/Logs/open-frontier-cloudflared.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
